# build/medusa.conf
[program:medusa]
user = nobody
command = /usr/sbin/python3 /opt/medusa/start.py --config /config/config.ini --datadir /config
umask = 000

[program:shutdown-script]
autorestart = false
startsecs = 0
user = root
command = /usr/local/bin/shutdown.sh 'medusa'
umask = 000
