# setup/medusa.conf
[program:medusa]
user = nobody
command = /usr/sbin/python2 /opt/medusa/start.py --config /config/config.ini --datadir /config
umask = 000
