# setup/medusa.conf
[program:sickrage]
user = nobody
command = /usr/sbin/python2 /opt/sickrage/SickBeard.py --config /config/config.ini --datadir /config